# scripts/com.typtel.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.typtel.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>BINARY_PATH</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>USER_HOME</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/typtel-menubar.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/typtel-menubar.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
